# deploy/notesos-frontend.service
[Unit]
Description=NotesOS Frontend (Next.js)
After=network.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_DIR__/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npx next start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
